# bt_volume_service_rpi0/bt_volume.service
[Unit]
Description=Bluetooth and Volume Control Service
After=multi-user.target bluetooth.service pulseaudio.service
Requires=bluetooth.service pulseaudio.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/bt_volume_service.py
Restart=on-failure
User=root
Group=root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
